# raspberry_pi_client.service
[Service]
ExecStart=/home/vladyslav/rc-car-repo/raspberry-pi-client/c/build/raspberrypiclient
Restart=always
RestartSec=5
User=root
